# broedenbot.service.example
[Unit]
Description=BroEdenBot
StartLimitIntervalSec=300
StartLimitBurst=3
After=network.target
RequiresMountsFor=/home/sadcatdad/BroEdenBot

[Service]
User=sadcatdad
WorkingDirectory=/home/sadcatdad/BroEdenBot
ExecStartPre=/usr/bin/test -x /home/sadcatdad/BroEdenBot/.venv/bin/python
ExecStartPre=/usr/bin/test -r /home/sadcatdad/BroEdenBot/data.db
ExecStartPre=/usr/bin/sqlite3 -readonly /home/sadcatdad/BroEdenBot/data.db "PRAGMA quick_check(1);"
ExecStart=/home/sadcatdad/BroEdenBot/.venv/bin/python /home/sadcatdad/BroEdenBot/main.py
TimeoutStartSec=120
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
