# broeden-dashboard.service.example
[Unit]
Description=BroEdenBot local dashboard
After=network.target

[Service]
Type=simple
User=sadcatdad
WorkingDirectory=/home/sadcatdad/BroEdenBot
EnvironmentFile=/home/sadcatdad/BroEdenBot/.env
ExecStart=/home/sadcatdad/BroEdenBot/.venv/bin/uvicorn dashboard.app:app --host 0.0.0.0 --port 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
